# wordle-bot.service
[Unit]
Description=Wordle Bot - Discord bot for tracking Wordle scores
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/wordle-bot
Environment=PATH=/home/ubuntu/wordle-bot/venv/bin
ExecStart=/home/ubuntu/wordle-bot/venv/bin/python main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target